# units/sync.target
[Unit]
Description=Mirror Sync Target
Requires=sync-daemon.service
After=multi-user.target
